# installer/rpm/algorand/0dnf-algorand-hourly.cron
#!/bin/bash
exec /usr/sbin/dnf-cron /etc/dnf/dnf-cron-algorand.conf
