# installer/rpm/algorand/0dnf-algorand-hourly.cron
#!/bin/bash
exec /usr/bin/dnf-automatic /etc/dnf/dnf-cron-algorand.conf
