# installer/rpm/algorand/0yum-algorand-hourly.cron
#!/bin/bash
exec /usr/sbin/yum-cron /etc/yum/yum-cron-algorand.conf
